# setup/files/systemd/photo-app.target
[Unit]
Description=Ready-to-run target for Photo App (started only when WiFi is up)
Requires=photo-app.service
