# setup/files/systemd/sync-photos.service
[Unit]
Description=Photo Library Sync (mirror cloud to live library)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=SYNC_TOOL=rclone
Environment="RCLONE_FLAGS=--checkers=8 --transfers=4 --retries=5 --copy-links"
Environment="RSYNC_FLAGS=-av --delete"
ExecStart=/bin/sh -lc '
set -euo pipefail
DEST="${PHOTO_LIBRARY_PATH:-${LIB_PATH:-}}"
if [[ -z "${DEST}" ]]; then
    echo "PHOTO_LIBRARY_PATH environment variable must be set" >&2
    exit 2
fi
case "${SYNC_TOOL}" in
    rsync)
        if [[ -z "${RSYNC_SOURCE:-}" ]]; then
            echo "RSYNC_SOURCE must be set when SYNC_TOOL=rsync" >&2
            exit 2
        fi
        exec rsync ${RSYNC_FLAGS:-"-av --delete"} "${RSYNC_SOURCE}" "${DEST}"
        ;;
    *)
        if [[ -z "${RCLONE_REMOTE:-}" ]]; then
            echo "RCLONE_REMOTE must be set when SYNC_TOOL=rclone" >&2
            exit 2
        fi
        exec rclone sync "${RCLONE_REMOTE}" "${DEST}" ${RCLONE_FLAGS:-"--checkers=8 --transfers=4 --retries=5 --copy-links"}
        ;;
esac'

[Install]
WantedBy=multi-user.target
