# setup/files/systemd/wifi-manager.service
[Unit]
Description=PhotoFrame WiFi Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@INSTALL_ROOT@/bin/wifi-manager watch --config @INSTALL_ROOT@/etc/wifi-manager.yaml
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_ROOT@/var
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
