# setup/files/systemd/wifi-watcher.service
[Unit]
Description=Photo Frame WiFi Watcher
After=network-manager.service NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
User=root
Group=root
Environment=RUST_LOG=info
Environment=FRAME_USER=frame
Environment=WIFI_IFNAME=wlan0
Environment=HOTSPOT_IP=192.168.4.1
Environment=PHOTO_FRAME_CONFIG=/etc/photo-frame/config.yaml
ExecStart=/opt/photo-frame/bin/wifi-watcher
Restart=on-failure

[Install]
WantedBy=multi-user.target
