# developer/systemd/wifi-overlay-test.service
[Unit]
Description=Wi‑Fi Overlay Takeover Test (dev)
After=greetd.service
Wants=greetd.service

[Service]
Type=oneshot
User=kiosk
Group=kiosk
ExecStart=/usr/local/bin/overlay-test
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
