# setup/system/cage@.service
# Managed by setup/system/install-systemd-units.sh
[Unit]
Description=Photo Frame Wayland session on %I
Documentation=https://github.com/vincentl/rust-photo-frame
After=systemd-user-sessions.service systemd-logind.service
After=plymouth-quit-wait.service
Before=graphical.target
Conflicts=getty@%i.service
Requires=systemd-logind.service
Wants=graphical.target

[Service]
Type=simple
User=kiosk
WorkingDirectory=/var/lib/photo-frame
Environment=PHOTOF_APP_BIN=/opt/photo-frame/bin/rust-photo-frame
Environment=PHOTOF_CONFIG=/opt/photo-frame/etc/config.yaml
ExecStart=/usr/bin/cage -- /opt/photo-frame/bin/rust-photo-frame --config /opt/photo-frame/etc/config.yaml
Restart=always
RestartSec=5s
TTYPath=/dev/%I
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty-fail
StandardOutput=journal
StandardError=journal
PAMName=cage
UtmpIdentifier=%I
UtmpMode=user

[Install]
Alias=photoframe-session@.service
Alias=display-manager.service
WantedBy=graphical.target
DefaultInstance=tty1
