# assets/systemd/buttond.service
# Managed by setup/system/modules/60-systemd.sh
[Unit]
Description=Photo Frame Power Button Monitor
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/var/lib/photoframe
ExecStart=/opt/photoframe/bin/buttond --config /etc/photoframe/config.yaml
# buttond derives /run/user/<uid> from its own (kiosk) UID at runtime; do NOT
# hardcode XDG_RUNTIME_DIR here. The kiosk account is created without a fixed
# UID, so a literal /run/user/1001 breaks on any host where kiosk != 1001.
Restart=always
RestartSec=1s
CapabilityBoundingSet=CAP_KILL CAP_SYS_BOOT
AmbientCapabilities=CAP_KILL CAP_SYS_BOOT
NoNewPrivileges=true

[Install]
WantedBy=graphical.target
