# assets/systemd/photoframe-sync.service
# Managed by setup/system/modules/60-systemd.sh
[Unit]
Description=Photo Frame Library Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kiosk
Group=kiosk
WorkingDirectory=/var/lib/photoframe
Environment=PHOTO_LIBRARY_PATH=/var/lib/photoframe/photos
Environment=SYNC_TOOL=rclone
Environment="RCLONE_FLAGS=--checkers=8 --transfers=4 --retries=5 --copy-links"
Environment="RSYNC_FLAGS=-av --delete"
EnvironmentFile=-/etc/photoframe/sync.env
# Syncs remote content into the cloud-backed mirror inside ${PHOTO_LIBRARY_PATH}.
ExecStart=/opt/photoframe/bin/photoframe-cloud-sync

[Install]
WantedBy=multi-user.target
