# assets/systemd/photoframe-wifi-manager.service
# Managed by setup/10-kiosk-bookworm.sh
[Unit]
Description=Photo Frame Wi-Fi Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/var/lib/photo-frame
Environment=RUST_LOG=info
ExecStart=/opt/photo-frame/bin/wifi-manager watch --config /opt/photo-frame/etc/wifi-manager.yaml
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
