# assets/systemd/photoframe-buttond.service
# Managed by setup/bootstrap/modules/60-systemd.sh
[Unit]
Description=Photo Frame Power Button Monitor
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/var/lib/photo-frame
ExecStart=/opt/photo-frame/bin/photo-buttond --single-window-ms 250 --double-window-ms 400 --debounce-ms 20 \
  --pidfile /run/photoframe/app.pid --procname rust-photo-frame --shutdown /opt/photo-frame/bin/photo-safe-shutdown
Restart=always
RestartSec=1s
CapabilityBoundingSet=CAP_KILL CAP_SYS_BOOT
AmbientCapabilities=CAP_KILL CAP_SYS_BOOT
NoNewPrivileges=true
RuntimeDirectory=photoframe

[Install]
WantedBy=multi-user.target
