# assets/systemd/photoframe-wifi-manager.service
# Managed by setup/system/modules/60-systemd.sh
[Unit]
Description=Photo Frame Wi-Fi Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kiosk
Group=kiosk
WorkingDirectory=/var/lib/photoframe
Environment=RUST_LOG=info
ExecStart=/opt/photoframe/bin/wifi-manager watch --config /opt/photoframe/etc/wifi-manager.yaml
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
