# setup/app/systemd/photo-frame.service
[Unit]
Description=Rust Photo Frame Kiosk
Wants=network-online.target seatd.service
After=network-online.target seatd.service

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
EnvironmentFile=-/etc/default/photo-frame
WorkingDirectory=@INSTALL_ROOT@/var
Environment=XDG_RUNTIME_DIR=%t/photo-frame
Environment=WAYLAND_DISPLAY=wayland-0
Environment=MOZ_ENABLE_WAYLAND=1
RuntimeDirectory=photo-frame
RuntimeDirectoryMode=0700
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
UtmpIdentifier=tty1
UtmpMode=login
ExecStart=@INSTALL_ROOT@/bin/photo-frame-kiosk
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
