# setup/modules/buttond/files/photo-buttond.service
[Unit]
Description=Photo Frame Power-Button Daemon (evdev)
After=multi-user.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/photo-frame/bin/photo-buttond --single-window-ms 250 --double-window-ms 400 \
  --debounce-ms 20 --pidfile /run/photo-app.pid --procname rust-photo-frame \
  --shutdown /opt/photo-frame/bin/photo-safe-shutdown
User=frame
Group=frame
Restart=always
RestartSec=1s
CapabilityBoundingSet=CAP_KILL CAP_SYS_BOOT
AmbientCapabilities=CAP_KILL CAP_SYS_BOOT
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
